# etc/service/confd/run
#!/bin/sh

exec /sbin/confd -interval=5
